# vagrant_files/prometheus_cl1_rf1/files/pushgateway.service
[Unit]
Description=Pushgateway

[Service]
ExecStart=/home/vagrant/pushgateway/pushgateway

[Install]
WantedBy=multi-user.target
